# com.dns-changer.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Identificador único do daemon -->
    <key>Label</key>
    <string>com.dns-changer.daemon</string>
    
    <!-- Programa a executar -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/dns_changer.py</string>
        <string>--interval</string>
        <string>300</string>
    </array>
    
    <!-- Carregar na inicialização -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Manter sempre rodando -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <!-- Intervalo mínimo entre reinicializações (em segundos) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Arquivo de log padrão -->
    <key>StandardOutPath</key>
    <string>__HOME__/.dns_changer/daemon.log</string>
    
    <!-- Arquivo de log de erros -->
    <key>StandardErrorPath</key>
    <string>__HOME__/.dns_changer/daemon_error.log</string>
    
    <!-- Variáveis de ambiente -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>
    
    <!-- Usuário que executará o daemon -->
    <key>UserName</key>
    <string>__USER__</string>
    
    <!-- Diretório de trabalho -->
    <key>WorkingDirectory</key>
    <string>__HOME__</string>
    
    <!-- Limpar variáveis de ambiente -->
    <key>ClearEnvironment</key>
    <false/>
    
    <!-- Modo de reinicialização -->
    <key>ExitTimeOut</key>
    <integer>20</integer>
    
    <!-- Prioridade de processo -->
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
